# pzserver-reboot.cron
# Nightly update + restart at 4am UTC (11pm PST / 2am EST)
58 3 * * * root /home/pzserver/weekly-reboot.sh
