# pzserver-reboot.cron
# Nightly update + restart with warnings at 3:58am
58 3 * * * root /home/pzserver/weekly-reboot.sh
