# pzserver.service
[Unit]
Description=Project Zomboid Dedicated Server
After=network-online.target
Wants=network-online.target

[Service]
Type=forking
User=pzserver
WorkingDirectory=/home/pzserver
ExecStart=/usr/bin/taskset -c 0,2,4,6 /home/pzserver/pzserver start
ExecStop=/home/pzserver/pzserver stop
Restart=on-failure
RestartSec=30
TimeoutStartSec=300
TimeoutStopSec=120

[Install]
WantedBy=multi-user.target
